# launchd/com.cortex.nightly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.cortex.nightly</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/gm/cortex/venv/bin/python</string>
        <string>/Users/gm/cortex/backend/digest/nightly.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/gm/cortex</string>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>2</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/gm/cortex/logs/nightly.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/gm/cortex/logs/nightly-error.log</string>
</dict>
</plist>
